# chatlol3-3b-agentic/build/src/supervisord.conf
[supervisord]
nodaemon=true

[program:ollama]
command=/usr/bin/ollama serve
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/ollama.log
stderr_logfile=/var/log/ollama_err.log

[program:python_main]
command=python3 /app/chatbot-ui/backend/main.py
autostart=true
autorestart=true
startretries=5
stdout_logfile=/var/log/python_backend.log
stderr_logfile=/var/log/python_backend_err.log
# Wait for Ollama before starting
startsecs=10

[program:chatbot_ui]
command=npm run dev -- --host
directory=/app/chatbot-ui
autostart=true
autorestart=true
stdout_logfile=/var/log/chatbot_ui.log
stderr_logfile=/var/log/chatbot_ui_err.log
